# deploy/sipfax.service.d/ppp.conf
# Drop-in to let the unprivileged sipfax service run pppd.
#
# pppd (from the Debian `ppp` package) is installed setuid-root, mode 4754
# root:dip, and needs CAP_NET_ADMIN to create the ppp0 interface. The base
# unit's hardening (NoNewPrivileges, RestrictSUIDSGID, a CAP_NET_BIND_SERVICE-
# only bounding set, plus the seccomp-based LockPersonality/SystemCallArchitectures
# which silently re-imply NoNewPrivileges) all block that. Relax exactly what
# pppd needs; the node process itself stays unprivileged and gains nothing
# (only the setuid pppd child escalates).
#
# Install: copy to /etc/systemd/system/sipfax.service.d/ppp.conf, then
#   systemctl daemon-reload && systemctl restart sipfax.service
#
# Prerequisites on the host:
#   - apt-get install ppp        (provides /usr/sbin/pppd, ppp_generic module)
#   - a kernel with PPP support   (the Debian "cloud" kernel has NONE; install
#                                  linux-image-amd64 and boot it)
#   - echo ppp_generic > /etc/modules-load.d/ppp.conf
#   - pre-create the secrets files the service rewrites in place:
#       touch /etc/ppp/chap-secrets /etc/ppp/pap-secrets
#       chown sipfax:sipfax /etc/ppp/chap-secrets /etc/ppp/pap-secrets
#       chmod 600 /etc/ppp/chap-secrets /etc/ppp/pap-secrets

[Service]
NoNewPrivileges=false
RestrictSUIDSGID=false
LockPersonality=false
SystemCallArchitectures=
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_ADMIN CAP_NET_RAW CAP_SETUID CAP_SETGID CAP_SYS_TTY_CONFIG CAP_DAC_OVERRIDE CAP_CHOWN CAP_FOWNER CAP_KILL CAP_AUDIT_WRITE
SupplementaryGroups=dip
ReadWritePaths=/run/lock /etc/ppp
